# com.annict.www/roles/redis/templates/redis-server.conf
description "redis server"

start on runlevel [23]
stop on shutdown

exec redis-server /etc/redis.conf

respawn
